# includes.liveiso/etc/xdg/autostart/installer.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Project Shards Installer
Exec=/usr/bin/shardsinstaller
Terminal=true
